# provision/kafka.conf
#!upstart
description "Kafka Server"

start on filesystem and net-device-up IFACE!=lo
stop on shutdown

exec /opt/kafka/bin/kafka-server-start.sh /opt/kafka/config/server.properties >> /var/log/kafka/kafka.log 2>&1
